# mavlink-router.service.in
[Unit]
Description=MAVLink Router

[Service]
Type=simple
ExecStart=@BINDIR@/mavlink-routerd --syslog
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
